# ops/linux/venue-pos-kiosk.service
[Unit]
Description=Venue POS kiosk (AppImage / Electron)
After=graphical-session.target network-online.target
PartOf=graphical-session.target
Wants=network-online.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=ELECTRON_IS_KIOSK=true
Environment=VENUE_POS_AGENT_ROOT=/opt/venue-pos/local-agent
Environment=VENUE_POS_FORCE_SETUP=1
ExecStart=/opt/venue-pos/pos/start-kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
